# target/linux/ar71xx/base-files/etc/uci-defaults/03_network-switchX-migration
#!/bin/sh
#
# Copyright (C) 2013 OpenWrt.org
#


exit 0
